# puppet/modules/homedir_cpan/files/profile.d/local-lib-perl.sh
# deployed by puppet - do not modify, your changes will be discarded.

# `perl -Mlocal::lib' prints out shell script suitable for use with eval,
# to set up the user's environment to use $HOME/perl5 as a prefix for
# perl modules.  It is not harmful if $HOME/perl5 doesn't exist yet.
#
# See: perldoc local::lib
#
eval $(perl -Mlocal::lib)
