# puppet/modules/baselayout/files/centos/profile.d/99homepath.sh
if test -d "$HOME/bin"; then
    PATH="$HOME/bin:$PATH"
    export PATH
fi

